# archives/filesystem/special/Edge-2L/etc/udev/rules.d/97-mpp-device.rules
KERNEL=="mpp_service" RUN+="/usr/local/bin/create-chromium-vda-vea-devices.sh"
