# archives/filesystem/special/Edge-2L/etc/udev/rules.d/98-fix-gdm-wayland.rules
SUBSYSTEM=="graphics", ACTION=="add", DEVPATH=="/devices/platform/display-subsystem/graphics/fb0", RUN+="/usr/local/bin/fix-gdm-wayland-hotplug.sh"
