# archives/filesystem/special/Edge2/etc/udev/rules.d/61-usbdevice.rules
SUBSYSTEM=="android_usb", ACTION=="change", RUN+="/usr/local/bin/usbdevice update"
SUBSYSTEM=="udc", RUN+="/usr/local/bin/usbdevice update"
