# archives/filesystem/common/etc/profile.d/adbd.sh
export ADB_TCP_PORT=5555
export ADBD_SHELL=/bin/bash
