# archives/filesystem/special/VIM-COMMON/etc/rc.local
#!/bin/sh -e
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "" on success or any other
# value on error.
#
# In order to enable or disable this script just change the execution
# bits.
#
# By default this script does nothing.

# Create system ssh keys if not exist
for a in /etc/ssh/ssh_host*; do
    [ -e "$a" ] && break
    dpkg-reconfigure -f noninteractive openssh-server
done

# For kodi
[ -f /sys/class/video/freerun_mode ] && chmod 777 /sys/class/video/freerun_mode
[ -f /sys/devices/system/cpu/cpu0/cpufreq/scaling_min_freq ] && chmod 777 /sys/devices/system/cpu/cpu0/cpufreq/scaling_min_freq
[ -f /sys/devices/system/cpu/cpu0/cpufreq/scaling_max_freq ] && chmod 777 /sys/devices/system/cpu/cpu0/cpufreq/scaling_max_freq
[ -f /sys/devices/system/cpu/cpu0/cpufreq/scaling_governor ] && chmod 777 /sys/devices/system/cpu/cpu0/cpufreq/scaling_governor
[ -f /sys/class/display/mode ] && chmod 777 /sys/class/display/mode
[ -f /sys/class/graphics/fb0/blank ] && chmod 777 /sys/class/graphics/fb0/blank
[ -f /sys/class/graphics/fb0/osd_clear ] && chmod 777 /sys/class/graphics/fb0/osd_clear

# For kplayer
[ -f /sys/class/tsync/vpause_flag ] && sudo chmod 777 /sys/class/tsync/vpause_flag
[ -f /sys/class/tsync/startsync_mode ] && sudo chmod 777 /sys/class/tsync/startsync_mode
[ -f /sys/class/tsync/slowsync_enable ] && sudo chmod 777 /sys/class/tsync/slowsync_enable
[ -f /sys/class/video/slowsync_repeat_enable ] && sudo chmod 777 /sys/class/video/slowsync_repeat_enable
[ -f /sys/class/video/vsync_slow_factor ] && sudo chmod 777 /sys/class/video/vsync_slow_factor

# For MCU
[ -f /sys/class/mcu/poweroff ] && chmod 777 /sys/class/mcu/poweroff
[ -f /sys/class/mcu/usb_pcie_switch_mode ] && chmod 777 /sys/class/mcu/usb_pcie_switch_mode

# For HDMI service
[ -f /sys/class/graphics/fb1/blank ] && chmod 777 /sys/class/graphics/fb1/blank
[ -f /sys/devices/virtual/graphics/fbcon/cursor_blink ] && chmod 777 /sys/devices/virtual/graphics/fbcon/cursor_blink
[ -f /sys/class/ppmgr/ppscaler ] && chmod 777 /sys/class/ppmgr/ppscaler

# For Kernel 4.9 FAN
[ -f /sys/class/fan/mode ] && chmod 777 /sys/class/fan/mode
[ -f /sys/class/fan/level ] && chmod 777 /sys/class/fan/level
[ -f /sys/class/fan/enable ] && chmod 777 /sys/class/fan/enable

# For Kernel mainline FAN
[ -f /sys/class/thermal/thermal_zone0/trip_point_3_temp ] && chmod 777 /sys/class/thermal/thermal_zone0/trip_point_3_temp

# For Freq-overclocking
[ -f /sys/devices/system/cpu/cpufreq/policy0/scaling_max_freq ] && chmod 777 /sys/devices/system/cpu/cpufreq/policy0/scaling_max_freq
[ -f /sys/devices/system/cpu/cpufreq/policy2/scaling_max_freq ] && chmod 777 /sys/devices/system/cpu/cpufreq/policy2/scaling_max_freq

# For WOL
[ -f /sys/class/wol/enable ] && chmod 0660 /sys/class/wol/enable

# Setup LED: heartbeat if bootup completed
[ -f /sys/class/leds/sys_led/trigger ] && echo heartbeat > /sys/class/leds/sys_led/trigger

# Setup cursor_blink
if ! which lightdm; then
	echo 1 > /sys/devices/virtual/graphics/fbcon/cursor_blink
else
	echo 0 > /sys/devices/virtual/graphics/fbcon/cursor_blink
fi

# Reboot test handle
if cat /proc/cmdline | grep -q reboot_test > /dev/null; then
	if ! systemctl is-active reboot-test-khadas | grep "^active$" > /dev/null; then
		systemctl start reboot-test-khadas
	fi
fi

exit 0
