# archives/filesystem/special/Edge1/etc/profile.d/gst.sh
export GST_GL_API=gles2
export GST_GL_PLATFORM=egl
