# archives/filesystem/common/lib/systemd/system/adb-khadas.service
[Unit]
Description=ADB setup.
DefaultDependencies=no
After=systemd-update-utmp-runlevel.service

[Service]
ExecStart=/usr/local/bin/adb-setup
RemainAfterExit=no

[Install]
WantedBy=sysinit.target
