# archives/filesystem/special/VIM4/etc/udev/rules.d/59-plug-event.rules
SUBSYSTEM=="amhdmitx",ACTION=="change",ENV{hdmitx_hpd}=="1",RUN+="/usr/local/bin/hdmitx_hpd_event.sh"
