# archives/filesystem/special/Edge1/lib/systemd/system/adb-helper-khadas.service
[Unit]
Description=ADB setup helper.
DefaultDependencies=no
After=systemd-update-utmp-runlevel.service

[Service]
ExecStart=/usr/local/bin/adbd
RemainAfterExit=no

[Install]
WantedBy=sysinit.target
