# archives/filesystem/special/VIM-COMMON/etc/udev/rules.d/60-display.rules
SUBSYSTEM=="extcon", ACTION=="change", ENV{NAME}=="hdmitx_extcon_hdmi", RUN+="/usr/local/bin/display_setup.sh $env{STATE}"
