# archives/filesystem/common/etc/update-motd.d/10-help-text
#!/bin/sh

printf "\n"
printf " * Website:        https://www.khadas.com\n"
printf " * Documentation:  https://docs.khadas.com\n"
printf " * Forum:          https://forum.khadas.com\n"
printf "\n"
