# archives/filesystem/special/VIM1S/lib/systemd/system/usb-gadget-khadas.service
[Unit]
Description=USB Gadget
DefaultDependencies=no
After=systemd-update-utmp-runlevel.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/usb-gadget
RemainAfterExit=true

[Install]
WantedBy=sysinit.target
