# archives/filesystem/special/Generic/etc/update-motd.d/69-khadas-setup
#!/bin/sh

## hyphop ##

echo
echo "TIPS: Try \`setup\` for tuning system"
echo
