# archives/filesystem/special/VIM1S/etc/udev/rules.d/59-plug-event.rules
SUBSYSTEM=="amhdmitx",ACTION=="change",ENV{hdmitx_hpd}=="1",RUN+="/usr/local/bin/hot_plug_event.sh"
